# rsl-booking.service
[Unit]
Description=RSL Booking Automation Service
After=network.target

[Service]
# user running the script
User=ai
Group=ai

# Set the working directory so .env is found
WorkingDirectory=/home/ai/github/rsl-booking

# Use the environment file
EnvironmentFile=/home/ai/github/rsl-booking/.env

# Use the python executable from your uv virtual environment
ExecStart=/home/ai/github/rsl-booking/.venv/bin/python /home/ai/github/rsl-booking/booker.py

# Do not restart after completion
Restart=no

# Environment variables (optional, if distinct from .env)
# Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target